# systemd-service-files/prodClient.service
[Unit]
Description=SENG302 Prod Client service

[Service]
Type=simple
Restart=always
RestartSec=10
User=gitlab-runner
WorkingDirectory=/home/gitlab-runner
ExecStart=/bin/bash runProdClient.sh

[Install]
WantedBy=multi-user.target
